# sysSentry-1.0.2/service/sentryCollector.service
[Unit]
Description = Collection module added for sysSentry and kernel lock-free collection
 
[Service]
ExecStart=/usr/bin/sentryCollector
ExecStop=/bin/kill $MAINPID
KillMode=process
Restart=on-failure
RestartSec=10s
 
[Install]
WantedBy = multi-user.target
